# rpm/moodle.cron
*/5 * * * * apache /usr/sbin/moodle-cron
